# feeds/packages/utils/scanbuttond/files/scanbuttond.init
#!/bin/sh /etc/rc.common

START=90
DAEMON=/usr/bin/scanbuttond

start() {
        service_start $DAEMON --quiet
}

stop() {
        service_stop $DAEMON --quiet
}

